# ops/systemd/schoolhub-backup.timer
[Unit]
Description=Run SchoolHub encrypted backup daily

[Timer]
OnCalendar=*-*-* 02:15:00 Asia/Jakarta
RandomizedDelaySec=15m
Persistent=true

[Install]
WantedBy=timers.target
