# ops/systemd/schoolhub-db-backup.timer
[Unit]
Description=Run SchoolHub PostgreSQL backup daily

[Timer]
OnCalendar=*-*-* 02:15:00
Persistent=true
Unit=schoolhub-db-backup.service

[Install]
WantedBy=timers.target
