# ops/systemd/schoolhub-smoke-monitor.timer
[Unit]
Description=Run SchoolHub smoke monitor every 10 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=10min
RandomizedDelaySec=30
Persistent=true
Unit=schoolhub-smoke-monitor.service

[Install]
WantedBy=timers.target
